# config/mirage-app.service
[Unit]
Description=The Mirage Community Next.js Application
Documentation=https://themirage.xxx/docs
After=network.target postgresql.service

[Service]
Type=exec
User=linuxuser
Group=linuxuser
WorkingDirectory=/home/linuxuser/mirage-community
Environment="NODE_ENV=production"
Environment="PORT=3000"
ExecStartPre=/usr/bin/npm run prisma:deploy
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/linuxuser/mirage-community
ReadWritePaths=/var/www/uploads
ReadWritePaths=/tmp

# Resource limits
LimitNOFILE=65536
LimitNPROC=512

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mirage-app

[Install]
WantedBy=multi-user.target 